# templates/profile.d/composer.sh
#!/bin/bash

export COMPOSER_HOME={{code_dir}}/.composer
