# templates/profile.d/npm.sh
# Add node_modules/.bin to the PATH
export PATH={{code_dir}}/node_modules/.bin:$PATH
